# Crux/2.1/supervisord.conf
[supervisord]
nodaemon=true

[program:crux]
autostart=true
directory=/crux-2.1.Linux.x86_64/bin/crux
command=/crux-2.1.Linux.x86_64/bin/crux


